# supervisord.conf
[program:cancel]
command=D:\\PYTHON\\python38\\python.exe D:\\code\\vsc\\dvid_Python\\FTP\\main.py
    
[supervisord]
nodaemon=true
 
[supervisorctl]